# CachyOS/config/systemd/user/backup.service
[Unit]
Description=Weekly Personal Backup
After=default.target

[Service]
Type=oneshot
ExecStart=%h/scripts/backup.sh
WorkingDirectory=%h
