# CachyOS/config/systemd/user/brightness-night.service
[Unit]
Description=Gece modu: %40 parlaklık

[Service]
ExecStart=%h/.local/bin/brightness-set.sh 40
