# CachyOS/config/systemd/user/brightness-day.service
[Unit]
Description=Gündüz modu: %80 parlaklık

[Service]
ExecStart=%h/.local/bin/brightness-set.sh 80
