# CachyOS/config/systemd/user/brightness-night.service
[Unit]
Description=Gece modu: %40 parlaklık

[Service]
Type=oneshot
ExecStart=%h/.local/bin/brightness-set.sh 40

[Install]
WantedBy=default.target
